# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/fail2ban/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# Runs the fail2ban daemon
# ==============================================================================
bashio::net.wait_for 10
if $(bashio::config 'fail2ban') -eq true; then
bashio::log.info "Starting Fail2ban...$(bashio::config 'fail2ban')"
exec fail2ban-server
fi
